# openvpn/files/openvpn.init
#!/bin/sh /etc/rc.common

. /opt/lib/functions.sh


OPENVPN="/opt/usr/sbin/openvpn"
SERVER_CONF="/opt/etc/config/server.conf"
CLIENT_CONF="/opt/etc/config/client.conf"

start(){
	service_start $OPENVPN --config  $SERVER_CONF &
}

stop(){
	service_stop $OPENVPN
}
