# openvpn/files/openvpn.init
#!/bin/sh /etc/rc.common

. /opt/lib/functions.sh

start()
{
	(/opt/usr/sbin/openvpn --config  /opt/etc/config/server)&
}


