# opkg-install/base-files/info_mipsel/dropbear.list
/opt/bin/dbclient
/opt/bin/scp
/opt/bin/ssh
/opt/sbin/dropbear
/opt/sbin/dropbearkey
/opt/sbin/dropbearconvert
/opt/etc/init.d/S01dropbear
/opt/etc/default/dropbear
