# xupnpd/files/xupnpd.init
#!/bin/sh /etc/rc.common

XUPNPD="/opt/usr/bin/xupnpd"

start(){
	service_start $XUPNPD
}

stop(){
	service_stop $XUPNPD
}

reload(){
	service_reload $XUPNPD
}
